# scheduler/crontab
# every 15 min: pull email, then normalize
*/15 * * * * cd /workspace && docker-compose run --rm ingestor-email && docker-compose run --rm normalizer

# hourly: rules
0 * * * * cd /workspace && docker-compose run --rm classifier

# daily 2:10am: budgets (idempotent)
10 2 * * * cd /workspace && docker-compose run --rm budgeter
